# autostart/mathpix-snipping-tool.desktop
[Desktop Entry]
Name=Mathpix Snipping Tool
Exec=''
Type=Application
Version=03.00.0104